# server/music-server.service
[Unit]
Description=Music Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/music-server/server
ExecStart=/root/music-server/server/music-server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PORT=80

[Install]
WantedBy=multi-user.target 